# expressvpn/3.17/crontab
SHELL=/bin/bash
*/1 * * * * /home/work/crontab/checkExpressVPNStatus.sh > /proc/1/fd/1 2>&1
